# systemd/uhdarc-backup.timer
[Unit]
Description=Run UHDARC Database Backup Daily

[Timer]
OnCalendar=*-*-* 02:00:00
Persistent=true

[Install]
WantedBy=timers.target
